# systemd/tailscale-monitor-collector.service
[Unit]
Description=Tailscale Network Monitor Collector
After=network.target
Documentation=https://github.com/yourusername/tailscale-monitor

[Service]
Type=simple
User=tailscale-monitor
Group=tailscale-monitor
WorkingDirectory=/opt/tailscale-monitor/collector
Environment="PATH=/opt/tailscale-monitor/collector/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/opt/tailscale-monitor"
Environment="PORT=48321"
ExecStart=/opt/tailscale-monitor/collector/venv/bin/python /opt/tailscale-monitor/collector/src/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tsmon-collector

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/tailscale-monitor /var/log/tailscale-monitor
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
